# atracker.service
[Unit]
Description=Activity Tracker Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/uv run --project %h/Projects/atracker atracker start
# To use the Rust binary instead, comment the line above and uncomment below:
# ExecStart=%h/Projects/atracker/rust/target/release/atracker start
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target
